# pwsh-profile/profile.d/99-banner.ps1
Write-Host "Use -> Show-Help <- para listar comandos do seu profile" -ForegroundColor Yellow
